# deploy/dashboard.service
[Unit]
Description=Crypto Trading AI Dashboard
After=network-online.target webapi.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/trading-ai/web
Environment=NODE_ENV=production
Environment=HOSTNAME=0.0.0.0
Environment=PORT=3000
ExecStart=/usr/bin/npm run start -- --hostname 0.0.0.0 --port 3000
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
