# deploy/trading-agent.service
[Unit]
Description=Crypto Trading AI — Orchestrator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trading-ai
Environment=PYTHONPATH=/home/ubuntu/trading-ai
EnvironmentFile=/home/ubuntu/trading-ai/.env
ExecStart=/home/ubuntu/trading-ai/venv/bin/python -m agent.orchestrator
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
